# systemd/surveybot-check.service
[Unit]
Description=ShopBot periodic health check
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/bots/shopbot
ExecStart=/opt/bots/shopbot/scripts/surveybot_check.sh
